# contrib/systemd/meow-chat.service
[Unit]
Description=MeowChat Server
After=network.target

[Service]
Type=simple
User=meow-chat
WorkingDirectory=/var/lib/meow-chat
EnvironmentFile=-/etc/meow-chat.env
ExecStart=/usr/local/bin/meow-chat-server
Restart=always
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
